# automation/launchd/com.niklas.finanztool.traderepublic-mail.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.niklas.finanztool.traderepublic-mail</string>
    <key>ProgramArguments</key>
    <array>
      <string>__NPM_PATH__</string>
      <string>run</string>
      <string>sync:traderepublic-mail</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__WORKING_DIRECTORY__</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>__PATH__</string>
      <key>FIREBASE_PROJECT_ID</key>
      <string>finanzperformance-tool</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>StandardOutPath</key>
    <string>/tmp/finanztool-traderepublic-mail.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/finanztool-traderepublic-mail.err.log</string>
  </dict>
</plist>
